# com.bse.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bse.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/pranjalsingh/nse project/venv/bin/python</string>
        <string>/Users/pranjalsingh/nse project/bse_monitor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/pranjalsingh/nse project</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/pranjalsingh/nse project/logs/bse-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/pranjalsingh/nse project/logs/bse-monitor-error.log</string>
</dict>
</plist>
